# package/sinzuo/udpcwmp/src/udpreport-init
#!/bin/sh /etc/rc.common
START=93
EXTRA_COMMANDS=""

start () {
        /usr/sbin/udpreport&
}

stop () {
        service_stop /usr/sbin/udpreport
}


